# ansible/templates/internet_of_agents.service.j2
# ansible/templates/internet_of_agents.service.j2
[Unit]
Description=Internet of Agents Service
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/internet_of_agents.env
ExecStart=/bin/bash -lc "source /etc/internet_of_agents.env && exec $VENV_PATH/bin/python /opt/internet_of_agents/agents/run_ui_agent_https.py --port 6000 --api-port 6001 --public-url http://$PUBLIC_IP:6000 --api-url https://$DOMAIN:6001"
Restart=always
User=root
WorkingDirectory=/opt/internet_of_agents

[Install]
WantedBy=multi-user.target