# ansible/templates/internet_of_agents.service.j2
# ansible/templates/internet_of_agents.service.j2
[Unit]
Description=Internet of Agents Service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
KillMode=none
EnvironmentFile=/etc/internet_of_agents.env
ExecStart=/bin/bash -lc "source /etc/internet_of_agents.env && /opt/internet_of_agents/agents/start_running_agents.sh"

[Install]
WantedBy=multi-user.target