# ansible/templates/internet_of_agents.service.j2
[Unit]
Description=Internet of Agents Service
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/internet_of_agents.env
ExecStart=/opt/internet_of_agents/venv/bin/bash -c "source /etc/internet_of_agents.env && /opt/internet_of_agents/agents/start_running_agents.sh"
Restart=on-failure

[Install]
WantedBy=multi-user.target