# ansible/templates/internet_of_agents.service.j2
# ansible/templates/internet_of_agents.service.j2
[Unit]
Description=Internet of Agents Service
After=network.target

[Service]
Type=simple
RemainAfterExit=always
KillMode=control-group
ExecStart=/bin/bash -lc "\
  source /etc/internet_of_agents.env && \
  python3 -u /opt/internet_of_agents/agents/run_ui_agent_https.py \
     --id agents${AGENT_ID_PREFIX} --port 6000 --api-port 6001 \
     --public-url http://$PUBLIC_IP:6000 --api-url https://$DOMAIN:6001 \
     --registry $REGISTRY_URL --ssl-cert $CERT_PATH --key $KEY_PATH & \
  exec python3 /opt/internet_of_agents/agents/agent_bridge.py"


[Install]
WantedBy=multi-user.target